# scripts/com.wireguard.wgctl-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  wgctl-agent — periodic reconciler for wg-mac devices joined to a
  Polar control plane. One agent per host; works for both hub and
  device roles (script branches on /etc/wgctl/config.json's "role").

  Install (handled by install.sh / wgctl-agent --install):
    cp wgctl-agent.sh                            /usr/local/sbin/wgctl-agent
    cp com.wireguard.wgctl-agent.plist           /Library/LaunchDaemons/
    launchctl enable    system/com.wireguard.wgctl-agent
    launchctl bootstrap system /Library/LaunchDaemons/com.wireguard.wgctl-agent.plist

  Disable:
    launchctl bootout   system/com.wireguard.wgctl-agent
    rm /Library/LaunchDaemons/com.wireguard.wgctl-agent.plist

  StartInterval = 60 → launchd relaunches the agent every minute. The
  agent is NOT a one-shot poller anymore: when there is a single mesh
  iface it runs an internal bounded long-poll loop (up to ~55 s,
  WGCTL_LP_BUDGET) that holds a GET /v1/peers?wait=… so a peer change
  propagates near-instantly. It auto-degrades to a single fetch against
  a server that doesn't support long-poll. LP_BUDGET (55) < StartInterval
  (60) so invocations never overlap; the 60 s relaunch is the coverage
  backstop. ExitTimeOut gives launchd headroom not to SIGKILL a held
  request on stop/reload.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wireguard.wgctl-agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/sbin/wgctl-agent</string>
    </array>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <!-- Headroom for an in-flight long-poll (budget ~55s + curl slack) so a
         stop/reload doesn't SIGKILL the held request mid-flight. -->
    <key>ExitTimeOut</key>
    <integer>75</integer>

    <key>StandardOutPath</key>
    <string>/var/log/wgctl-agent.out.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/wgctl-agent.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
